# robot2/rovy2.service
[Unit]
Description=ROVY Robot2 Service (for Cloud2)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rovy
Group=rovy
SupplementaryGroups=video plugdev audio
WorkingDirectory=/home/rovy/rovy_client/robot2
ExecStart=/bin/bash /home/rovy/rovy_client/robot2/start.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=CAMERA_WEBCAM_DEVICE=/dev/video0
Environment=ROVER_SERIAL_DEVICE=/dev/ttyAMA0
Environment=AUDIODEV=hw:3,0
Environment=XDG_RUNTIME_DIR=/run/user/1000

[Install]
WantedBy=multi-user.target

